# services/com.kyrelldixon.agent-platform.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kyrelldixon.agent-platform</string>
    <key>UserName</key>
    <string>kyrelldixon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kyrelldixon/.bun/bin/bun</string>
        <string>run</string>
        <string>server/index.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kyrelldixon/projects/agent-platform</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>3400</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kyrelldixon/Library/Logs/agent-platform.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kyrelldixon/Library/Logs/agent-platform.err</string>
</dict>
</plist>
